# scripts/claude-migrate.service
[Unit]
Description=Daily incremental backup of a Claude.ai profile
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Replace `source` with your profile name from `claude-migrate accounts`.
ExecStart=/usr/local/bin/claude-migrate --quiet backup source
StandardOutput=append:%h/.local/state/claude-migrate/backup.log
StandardError=append:%h/.local/state/claude-migrate/backup.log

[Install]
WantedBy=default.target
