# scripts/com.user.claudemigrate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.claudemigrate</string>
  <key>ProgramArguments</key>
  <!-- Replace `source` with your profile name from `claude-migrate accounts`. -->
  <array>
    <string>/usr/local/bin/claude-migrate</string>
    <string>--quiet</string>
    <string>backup</string>
    <string>source</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>17</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>~/Library/Logs/claude-migrate.log</string>
  <key>StandardErrorPath</key>
  <string>~/Library/Logs/claude-migrate.log</string>
</dict>
</plist>
